# testing/docker/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
username = magento
password = password

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
environment=SYMFONY_ENV=prod


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:cron]
directory=/var/www/magento
command=/bin/sh -c "while [ true ]; do ( php -dmemory_limit=-1 bin/magento cron:run ); sleep 60; done"
stdout_logfile=/dev/stdout
autostart=true
autorestart=true

[program:php]
command=php-fpm -R
autostart=true
autorestart=true